# infra/ansible/production/supervisord.conf
[supervisord]
logfile=/home/ubuntu/arb/logs/supervisord.log
childlogdir=/home/ubuntu/arb/logs
pidfile=/home/ubuntu/arb/pids/supervisord.pid

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[unix_http_server]
file = /tmp/supervisor.sock
chmod = 0777

[inet_http_server]
port = 127.0.0.1:9001

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

# [program:uwsgi_arb]
# command = /bin/bash /home/ubuntu/arb/deploy/uwsgi_start.sh
# startsecs = 10
# autostart = true
# autorestart = true
# killasgroup = true
# stopasgroup = true
# stdout_logfile = /home/ubuntu/arb/logs/supervisord_uwsgi_arb_stdout.log
# stdout_logfile_maxbytes=1MB
# stdout_logfile_backups=10
# stdout_capture_maxbytes=1MB
# stdout_events_enabled=false
# stderr_logfile = /home/ubuntu/arb/logs/supervisord_uwsgi_arb_stderr.log
# stderr_logfile_maxbytes=1MB
# stderr_logfile_backups=10
# stderr_capture_maxbytes=1MB
# stderr_events_enabled=false


# [program:stupid]
# # command = /usr/bin/python -m arb.app.tasks.stupid
# command = /usr/bin/python -m arb.app.tasks.stupid
# stdout_logfile = /home/ubuntu/arb/logs/SUPERVISORD_stupid_stdout.log
# stderr_logfile = /home/ubuntu/arb/logs/SUPERVISORD_stupid_stderr.log
# autostart = true
# killasgroup = true
# stopasgroup = true
# stdout_logfile_maxbytes=10MB
# stderr_logfile_maxbytes=10MB


[program:data_collection]
command = /bin/bash /home/ubuntu/arb/deploy/tasks/data_collection.sh
stdout_logfile = /home/ubuntu/arb/logs/SUPERVISORD_data_collection.log
stderr_logfile = /home/ubuntu/arb/logs/SUPERVISORD_data_collection.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
autostart = true
killasgroup = true
stopasgroup = true


//[program:mock_live]
//command = /bin/bash /home/ubuntu/arb/deploy/tasks/mock_live.sh
//stdout_logfile = /home/ubuntu/arb/logs/SUPERVISORD_mock_live.log
//stderr_logfile = /home/ubuntu/arb/logs/SUPERVISORD_mock_live.log
//stdout_logfile_maxbytes=10MB
//stderr_logfile_maxbytes=10MB
//autostart = true
//killasgroup = true
//stopasgroup = true


[program:live_trade_01]
command = /bin/bash /home/ubuntu/arb/deploy/tasks/live_trade_01.sh
stdout_logfile = /home/ubuntu/arb/logs/SUPERVISORD_live_trade_01.log
stderr_logfile = /home/ubuntu/arb/logs/SUPERVISORD_live_trade_01.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
autostart = true
killasgroup = true
stopasgroup = true
